# src/main/conf/filemgr/bin/filemgr
#!/bin/sh
#
# init script for XmlRpcFileManager
#
# chkconfig: 345 88 22
# description: CAS File Manager
#
# Copyright (c) 2007 California Institute of Technology.
# ALL RIGHTS RESERVED. U.S. Government Sponsorship acknowledged.
#
# $Id$

SERVER_PORT=9000
export SERVER_PORT
CMD=$1
set CMD
JAVA=`which java`
export JAVA
DIR=`dirname $0`
export DIR
cd $DIR
DIR_PATH=`pwd`
set DIR_PATH
CAS_FILEMGR_HOME=$DIR_PATH/../
export CAS_FILEMGR_HOME
RUN_HOME=${CAS_FILEMGR_HOME}/../run
export RUN_HOME
CAS_FILEMGR_PROPS=${CAS_FILEMGR_HOME}/etc/filemgr.properties
export CAS_FILEMGR_PROPS

# See how we were called.
case "$CMD" in
  start)
        echo -n "Starting cas file manager: "
        $JAVA -Xms256m -Xmx1024m -Djava.ext.dirs=${CAS_FILEMGR_HOME}/lib \
        -Djava.util.logging.config.file=${CAS_FILEMGR_HOME}/etc/logging.properties \
        -Dgov.nasa.jpl.oodt.cas.filemgr.properties=${CAS_FILEMGR_PROPS} \
        gov.nasa.jpl.oodt.cas.filemgr.system.XmlRpcFileManager --portNum $SERVER_PORT &       
        echo $! >${RUN_HOME}/cas.filemgr.pid 
        echo "OK"
        sleep 5
        ;;
  stop)
        echo -n "Shutting down cas file manager: "
        kill `cat ${RUN_HOME}/cas.filemgr.pid`
        echo "OK"
        ;;
  restart)
        $0 stop
        $0 start
        ;;
  *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
esac

exit 0
